# configs/systemd/backup.timer
[Unit]
Description=Dagelijkse back-up om 02:00

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target

