# configs/systemd/backup.service
[Unit]
Description=Back-up service (triggered by backup.timer)

[Service]
Type=oneshot
ExecStart=/opt/scripts/backup-cron.sh
User=root

