# configs/systemd/example-monitor.service
[Unit]
Description=Voorbeeld monitoring service
After=network.target

[Service]
Type=simple
User=salah
ExecStart=/opt/scripts/monitor.sh
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

